# home/dot_config/systemd/user/unclutter.service
[Unit]
Description=Hide cursor when not in use
After=graphical-session.target

[Service]
Type=exec
ExecStart=unclutter
Restart=always
RestartSec=2
Slice=background-graphical.slice

[Install]
WantedBy=graphical-session.target
